# arne/Corona/crontab
0 0,6,12,18 * * *  /home/arnee/git/bioinfo-toolbox/arne/Corona/cron_run
0 0 * * *  /home/arnee/bin/run_backup.bash
